# systemd/tesserae-api-backup.service
[Unit]
Description=Back up the Tesserae PostgreSQL database
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/opt/tesserae-api/backup_db.sh
